# local.d/keycodes.start
#!/bin/zsh
#
# Set keycodes for brightness up/down
#
###

setkeycodes e012 224
setkeycodes e017 225
