# local.d/pm.start
# Set radeon power management to "mid" profile
echo mid > /sys/class/drm/card0/device/power_profile

if [ `cat /sys/class/power_supply/AC/online` = "0" ]; then
	pm-powersave true
else
	pm-powersave false
fi
