# local.d/pm.start
if [ `cat /sys/class/power_supply/AC/online` = "0" ]; then
	SATA_ALPM_ENABLE=true /usr/sbin/pm-powersave true
else
	/usr/sbin/pm-powersave false
fi
