# local.d/radeon_pm.start
# Set radeon power management to "mid" profile
echo mid > /sys/class/drm/card0/device/power_profile
